# containers/overlay.d/10-vm-swapfile/usr/lib/systemd/system/var-swapfile.service
[Unit]
Description=Create /var/swapfile sized from available RAM
ConditionPathExists=!/var/swapfile
Before=var-swapfile.swap

[Service]
Type=oneshot
ExecStart=/bin/sh /usr/lib/systemd/scripts/var-swapfile-create.sh
RemainAfterExit=yes
